# sms-forwarder.service
[Unit]
Description=SMS Forwarder - Air724UG SMS Gateway
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/sms-forwarder
ExecStart=/usr/bin/node src/app.js
Restart=always
RestartSec=5
# 60 秒内最多重启 10 次，防止频繁重启
StartLimitIntervalSec=60
StartLimitBurst=10
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
